# scripts/crontab
MAILTO=mmckenna.phila@gmail.com

15 0,18 * * * root certbot renew